# Server App/systemd/MCPummelchenModServer.service.d/minecraft-autostart.conf
[Service]
Environment=PUMMELCHEN_MINECRAFT_AUTOSTART=true
Environment=PUMMELCHEN_MINECRAFT_DIR=/opt/pummelchen-swift/runtime/minecraft
Environment="PUMMELCHEN_MINECRAFT_START_COMMAND=./run.sh nogui"
Environment=PUMMELCHEN_MINECRAFT_HOST=127.0.0.1
Environment=PUMMELCHEN_MINECRAFT_PORT=25565
Environment=PUMMELCHEN_MINECRAFT_LOG=/opt/pummelchen-swift/runtime/logs/minecraft-live.log
ExecStartPre=/bin/sh -c "/usr/sbin/iptables -C INPUT -p tcp --dport 25575 ! -s 127.0.0.1 -j DROP 2>/dev/null || /usr/sbin/iptables -I INPUT 1 -p tcp --dport 25575 ! -s 127.0.0.1 -j DROP"
ExecStartPre=/bin/sh -c "/usr/sbin/ip6tables -C INPUT -p tcp --dport 25575 -j DROP 2>/dev/null || /usr/sbin/ip6tables -I INPUT 1 -p tcp --dport 25575 -j DROP"
