# Server App/systemd/MCPummelchenModServer.service.d/minecraft-autostart.conf
[Service]
Environment=PUMMELCHEN_MINECRAFT_AUTOSTART=true
Environment=PUMMELCHEN_MINECRAFT_DIR=/opt/pummelchen-swift/runtime/minecraft
Environment="PUMMELCHEN_MINECRAFT_START_COMMAND=./run.sh nogui"
Environment=PUMMELCHEN_MINECRAFT_HOST=127.0.0.1
Environment=PUMMELCHEN_MINECRAFT_PORT=25565
Environment=PUMMELCHEN_MINECRAFT_LOG=/opt/pummelchen-swift/runtime/logs/minecraft-live.log
Environment=PUMMELCHEN_MINECRAFT_RCON_FIREWALL=true
