# systemd/pummelchen-minecraft-metrics.service
[Unit]
Description=Pummelchen Minecraft Prometheus exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/minecraft_mods
ExecStart=/usr/bin/python3 /var/minecraft_mods/scripts/minecraft_metrics_exporter.py --listen-host 127.0.0.1 --listen-port 7792 --rcon-password-file /var/minecraft_mods/secrets/rcon.password
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=full
ReadWritePaths=/var/minecraft_mods/site /var/minecraft_mods/data

[Install]
WantedBy=multi-user.target
