# systemd/pummelchen-server.service
[Unit]
Description=Pummelchen Swift read-only API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/minecraft_mods
ExecStart=/var/minecraft_mods/swift/PummelchenSwift/.build/release/pummelchen-server serve --project-root /var/minecraft_mods --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=5
UMask=0022
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/minecraft_mods
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
