# systemd/pummelchen-live-stats.timer
[Unit]
Description=Update Pummelchen live stats JSON every 10 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true
Unit=pummelchen-live-stats.service

[Install]
WantedBy=timers.target
