# Server App/systemd/MCPummelchenModUpdateScan.service
[Unit]
Description=MCPummelchen daily Minecraft mod update scan
After=network-online.target MCPummelchenModServer.service
Wants=network-online.target
Conflicts=MCPummelchenModServer.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/minecraftai/26.1.2/runtime
EnvironmentFile=/var/minecraftai/26.1.2/etc/server.env
ExecStartPre=/usr/bin/systemctl stop MCPummelchenModServer.service
ExecStart=/usr/bin/flock -n /run/MCPummelchenModUpdateScan.lock /var/minecraftai/26.1.2/bin/MCPummelchenModServer mod-update-scan --project-root /var/minecraftai/26.1.2/runtime --duckdb /var/minecraftai/26.1.2/runtime/data/pummelchen.duckdb --all-supported true --seed-from-project-data true --discover-source-links true --discovery-searches-per-second 2 --max-urls-per-window 5 --window-seconds 10
ExecStopPost=/usr/bin/systemctl start MCPummelchenModServer.service
StandardOutput=journal
StandardError=journal
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/minecraftai/26.1.2/runtime /run
